# packages/configure_vpn/configure_vpn/configure-vpn.service
[Unit]
Description=Configuration Service for Initializing the CLIP VPN variant
Requires=ipsec.service
After=ipsec.service

[Service]
Type=oneshot
ExecStart=/bin/sh /etc/init.d/configure-vpn start
ExecStop=/bin/sh /etc/init.d/configure-vpn stop

[Install]
RequiredBy=ipsec.service
